# .config/systemd/user/1password.service
[Unit]
Description=1password

[Service]
ExecStartPre=/bin/sleep 3
ExecStart=1password

[Install]
WantedBy=autostart.target
