# .config/systemd/user/1password.service
[Unit]
Description=1Password

[Service]
Type=simple
ExecStart=/opt/1Password/1password
restart=on-failure

[Install]
WantedBy=default.target
