# .config/systemd/user/autostart.target
[Unit]
Description=Autostart Target
Documentation=man:systemd.special(7)

[Install]
WantedBy=default.target
