# .config/systemd/user/caffeine.service
[Unit]
Description=Caffeine
After=graphical.target

[Service]
ExecStart=/usr/bin/caffeine

[Install]
WantedBy=autostart.target
