# .config/systemd/user/caffeine.service
[Unit]
Description=Caffeine
After=graphical.target

[Service]
ExecStart=caffeine-indicator

[Install]
WantedBy=autostart.target
