# .config/systemd/user/clipman.service
[Unit]
Description=Clipman
After=graphical.target

[Service]
ExecStart=xfce4-clipman

[Install]
WantedBy=autostart.target
