# .config/systemd/user/darkman.service
[Unit]
Description=Darkman
After=graphical.target

[Service]
ExecStart=/usr/bin/darkman set dark

[Install]
WantedBy=autostart.target
