# .config/systemd/user/flameshot.service
[Unit]
Description=flameshot
After=graphical.target

[Service]
ExecStart=flameshot

[Install]
WantedBy=autostart.target
