# .config/systemd/user/google-chat.service
[Unit]
Description=Google Chat

[Service]
ExecStart=/opt/google-chat-linux/google-chat-linux %U

[Install]
WantedBy=autostart.target
