# .config/systemd/user/morgend.timer
[Unit]
Description=Morning

[Timer]
OnBootSec=30s

[Install]
WantedBy=timers.target
