# .config/systemd/user/morgend.service
[Unit]
Description=Morning

[Service]
Type=oneshot
ExecStart=notify-send "Let's goo!"

[Install]
WantedBy=autostart.target
