# .config/systemd/user/nextcloud.service
[Unit]
Description=nextcloud

[Service]
ExecStart=nextcloud

[Install]
WantedBy=autostart.target
