# .config/systemd/user/nextcloud.service
[Unit]
Description=nextcloud
After=graphical.target

[Service]
ExecStart=nextcloud

[Install]
WantedBy=autostart.target
