# .config/systemd/user/openweather.service
[Unit]
Description=Openweather Generate

[Service]
ExecStart=/home/stevendejong/.config/ntfd/generate.sh

[Install]
WantedBy=autostart.target
