# .config/systemd/user/slack.service
[Unit]
Description=Slack
After=graphical.target

[Service]
ExecStart=/usr/bin/slack -u

[Install]
WantedBy=autostart.target
