# .config/systemd/user/solaar.service
[Unit]
Description=Solaar
After=graphical.target

[Service]
ExecStart=solaar -w hide

[Install]
WantedBy=autostart.target
