# .config/systemd/user/whatsapp.service
[Unit]
Description=WhatsApp

[Service]
ExecStart=whatsapp-nativefier

[Install]
WantedBy=autostart.target
