# .config/systemd/user/whatsapp.service
[Unit]
Description=WhatsApp
After=graphical.target

[Service]
ExecStart=/usr/bin/whatsapp-for-linux

[Install]
WantedBy=autostart.target
