# service/autostart/jetbrains-toolbox.service
[Unit]
Description=JetBrains Toolbox

[Service]
ExecStart=/home/stevendejong/.local/share/JetBrains/Toolbox/bin/jetbrains-toolbox --minimize
