# service/autostart/owncloud.service
[Unit]
Description=flameshot

[Service]
ExecStart=flameshot
