# service/autostart/power-manager.service
[Unit]
Description=Power Manager

[Service]
ExecStart=xfce4-power-manager
