# service/autostart/slimbookbattery.service
[Unit]
Description=Slimbook Battery

[Service]
ExecStart=slimbookbattery
