# service/autostart/steam.service
[Unit]
Description=Steam

[Service]
ExecStart=/usr/bin/steam-runtime %U
