# service/autostart/xpad.service
[Unit]
Description=Xpad

[Service]
ExecStart=xpad
