# config/systemd/user/flameshot.service
[Unit]
Description=flameshot
After=bspwm.target

[Service]
ExecStart=flameshot

[Install]
WantedBy=autostart.target
