# config/systemd/user/jetbrains-toolbox.service
[Unit]
Description=JetBrains Toolbox
Requires=network-online.target
After=network-online.target

[Service]
ExecStartPre=/usr/bin/sleep 10
ExecStart=/home/stevendejong/.local/share/JetBrains/Toolbox/bin/jetbrains-toolbox --minimize

[Install]
WantedBy=autostart.target
