# config/systemd/user/nm-applet.service
[Unit]
Description=Network Manager Applet
After=dwm.target

[Service]
ExecStart=nm-applet

[Install]
WantedBy=autostart.target
