# config/systemd/user/power-manager.service
[Unit]
Description=Power Manager
After=dwm.target

[Service]
ExecStart=xfce4-power-manager

[Install]
WantedBy=autostart.target
