# config/systemd/user/slimbookbattery.service
[Unit]
Description=Slimbook Battery
After=bspwm.target

[Service]
ExecStart=python3 /usr/share/slimbookbattery/src/slimbookbatteryindicator.py

[Install]
WantedBy=autostart.target
