# config/systemd/user/xpad.service
[Unit]
Description=Xpad
After=bspwm.target

[Service]
ExecStart=xpad

[Install]
WantedBy=autostart.target
