# config/systemd/user/xpad.service
[Unit]
Description=Xpad
After=dwm.target

[Service]
ExecStart=xpad

[Install]
WantedBy=autostart.target
