# service/battery-notify/battery-notify.timer
[Unit]
Description=Check battery status every few minutes to warn the user in case of low battery
# Set `Requires` to the name of the battery alert service
Requires=battery-notify.service

# Define when and how the timer activates
[Timer]
# Start 1 minute after boot...
OnBootSec=1m
# ...and again every 5 minutes after `battery-notify.service` runs
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
