# roles/services/files/autostart/clipman.service
[Unit]
Description=Clipman

[Service]
ExecStart=xfce4-clipman
