# roles/services/files/autostart/kdocker.service
[Unit]
Description=Kdocker

[Service]
ExecStart=kdocker
