# roles/services/files/autostart/nm-applet.service
[Unit]
Description=Network Manager Applet

[Service]
ExecStart=nm-applet
