# roles/services/files/autostart/obsidian.service
[Unit]
Description=Obsidian

[Service]
ExecStart=obsidian