# roles/services/files/autostart/owncloud.service
[Unit]
Description=Owncloud Client

[Service]
ExecStart=owncloud
