# roles/services/files/autostart/slimbookbattery.service
[Unit]
Description=Slimbook Battery

[Service]
ExecStart=python3 /usr/share/slimbookbattery/src/slimbookbatteryindicator.py
